# installation/service/bell_sys_monitor_websocket.service
[Unit]
Description=Bell Sys Monitor WebSocket Server
After=network.target

[Service]
User=root
WorkingDirectory=/path/to/your/project
ExecStart=/usr/bin/python3 /path/to/your/project/app.py
Restart=always

[Install]
WantedBy=multi-user.target
